# systemd/maxwell-boot.service
[Unit]
Description=Maxwell set to quickly seed /dev/random

[Service]
Type=oneshot
ExecStart=/usr/bin/maxwell

[Install]
WantedBy=sysinit.target
